# edulinux/trash/polkit.conf
[program:polkit]
command=/usr/libexec/polkitd --no-debug
autostart=true
autorestart=true
stderr_logfile=/var/log/polkit.err.log
stdout_logfile=/var/log/polkit.out.log
user=root
