# pkg/xrdp/context/dbus.conf
[program:dbus]
command=dbus-daemon --system --fork
autostart=true
autorestart=true
stderr_logfile=/var/log/dbus.err.log
stdout_logfile=/var/log/dbus.out.log
