# pkg/apache2/trash/services.d/apache/finish
#!/bin/bash
kill -TERM $(cat /var/run/apache2/apache2.pid)
